# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:ms-inspect]
command=ms-inspect
environment=RADIO_MCP_TRANSPORT=http,RADIO_MCP_PORT=8000
autorestart=true
stdout_logfile=/var/log/ms-inspect.log
stderr_logfile=/var/log/ms-inspect.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:wildcat]
command=wildcat
environment=WILDCAT_MS_PATH=%(ENV_WILDCAT_MS_PATH)s
autorestart=unexpected
exitcodes=0
stdout_logfile=/var/log/wildcat.log
stderr_logfile=/var/log/wildcat.log
